# hello_world.conf
# hello_world - The hello world web app
description "Hello World Web App"

start on runlevel [2345]
stop on runlevel [!2345]

exec ruby /usr/local/hello_world/hello_world.rb